# monitor/antilles-confluent-mond/etc/antilles-confluent-mond-health.service
[Unit]
Description=Monitor health metric from confluent

[Service]
Type=oneshot
User=antilles
ExecStart=/usr/bin/antilles-confluent-mond health
