# monitor/antilles-confluent-mond/etc/antilles-confluent-mond-power.service
[Unit]
Description=Monitor power metric from confluent

[Service]
Type=oneshot
User=antilles
ExecStart=/usr/bin/antilles-confluent-mond power
