# monitor/antilles-confluent-mond/etc/antilles-confluent-mond.service
[Unit]
Description=Confluent monitor component
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "echo Welcome to use %n"

[Install]
WantedBy=multi-user.target
