# antilles-core/etc/supervisor.d/antilles.ini
[program:antilles]
command=/usr/sbin/antilles-gunicorn --paste /etc/antilles/antilles.ini -w 5 --threads 10 --bind :18080 --log-config /etc/antilles//antilles.ini --timeout 3600 --keep-alive 65 --log-level info --access-logfile - --error-logfile - --capture-output --preload
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/antilles/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10

[program:antilles_beat]
command=/usr/bin/antilles beat --log-level INFO --conf-path /var/run/antilles/core/beats
user=antilles
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/antilles/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10

[program:antilles_worker]
command=/usr/bin/antilles worker --log-level INFO --concurrency 15
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/antilles/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10

[program:antilles_alarm_creator]
command=/usr/bin/antilles alarm_creator --log-level DEBUG
user=antilles
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/antilles/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
