# antilles-core/etc/supervisor.d/jobdaemon.ini
[program:jobdaemon]
command=/usr/bin/python2.7 /opt/antilles/core/job_daemon_project/listener
redirect_stderr=true
stdout_logfile=/var/log/antilles/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10