# monitor/antilles-vnc-mond/etc/antilles-vnc-mond-xvnc.service
[Unit]
Description=Monitor Xvnc process

[Service]
Type=oneshot
ExecStart=/usr/bin/antilles-vnc-mond
